# security-monitor.service
[Unit]
Description=Security Monitor Kyiv (Flask App)
After=network.target

[Service]
User=root
WorkingDirectory=/root/geminicli/security-monitor-kyiv
ExecStart=/root/geminicli/security-monitor-kyiv/venv/bin/gunicorn -b 0.0.0.0:5000 app:app
Restart=always

[Install]
WantedBy=multi-user.target
